# supervisor_conf.ini
[program:dq]
command =/home/oyo_admin/miniconda3/envs/data-quality/bin/uwsgi --ini /data/bigdata-dq/uwsgi.ini
autostart = true
stopsignal=QUIT    
user=root
startsecs = 5 
startretries = 3 
autorestart = true 
redirect_stderr = true 
stdout_logfile_maxbytes = 20MB 
stdout_logfile = /data/logs/supervisor/dq_stdout.log
stderr_logfile = /data/logs/supervisor/dq_err.log
loglevel=info
